# claude-irc-bot.service.example
[Unit]
Description=Claude IRC Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/claude-irc-bot
ExecStart=/usr/local/bin/node src/index.js
ProtectSystem=strict
ReadWritePaths=/home/YOUR_USER
PrivateTmp=true
NoNewPrivileges=true
Restart=always
RestartSec=10
EnvironmentFile=/home/YOUR_USER/claude-irc-bot/.env

[Install]
WantedBy=multi-user.target
